# systemd/monitors-revert.service
[Unit]
Description=Revert auf Dreier-Monitorlayout
[Service]
Type=oneshot
ExecStart=%h/.local/bin/monitors-mode.sh triple --primary DP-3
